# supervisord.conf
[program:node]
command = node all/server.js
directory = /home/dotcloud/current
